# hedera-mirror-rosetta/build/supervisord-offline.conf
; Sample supervisor config file.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html

[supervisord]
nodaemon=true               ; start in foreground if true; default false

[unix_http_server]
file=/tmp/supervisor.sock

[program:rosetta]
command=/app/rosetta/hedera-mirror-rosetta
user=rosetta
directory=/app/rosetta/
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stopsignal=SIGTERM
environment=HEDERA_MIRROR_ROSETTA_ONLINE=false
