# hedera-mirror-rosetta/build/supervisord.conf
; Sample supervisor config file.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html

[supervisord]
nodaemon=true               ; start in foreground if true; default false

[unix_http_server]
file=/tmp/supervisor.sock

[program:postgres]
command=/usr/lib/postgresql/9.6/bin/postgres -c config_file=/etc/postgresql/9.6/main/postgresql.conf
user=postgres
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
priority=1

[program:importer]
command=java -Xms2g -Xmx10g -XX:+CrashOnOutOfMemoryError -DLog4jContextSelector=org.apache.logging.log4j.core.async.AsyncLoggerContextSelector -jar /var/importer/hedera-mirror-importer.jar --spring.config.additional-location=file:/var/importer/
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
priority=10

[program:rosetta]
command=/var/rosetta/rosetta-executable
directory=/var/rosetta/
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stopsignal=SIGTERM
priority=100
