# hedera-mirror-rosetta/build/supervisord.conf
; Sample supervisor config file.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html

[supervisord]
nodaemon=true               ; start in foreground if true; default false

[unix_http_server]
file=/tmp/supervisor.sock

[program:postgres]
command=/usr/lib/postgresql/%(ENV_PG_VERSION)s/bin/postgres -c config_file=%(ENV_PGCONF)s/postgresql.conf
user=postgres
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
priority=1

[program:importer]
command=java -Xms512m -Xmx4g -XX:+CrashOnOutOfMemoryError -DLog4jContextSelector=org.apache.logging.log4j.core.async.AsyncLoggerContextSelector -jar /app/importer/hedera-mirror-importer.jar --spring.config.additional-location=file:/app/importer/
user=importer
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
priority=10
environment=HEDERA_MIRROR_IMPORTER_PARSER_RECORD_ENTITY_PERSIST_NONFEETRANSFERS=true,HEDERA_MIRROR_IMPORTER_PARSER_RECORD_ENTITY_REDIS_ENABLED=false

[program:rosetta]
command=/app/rosetta/hedera-mirror-rosetta
user=rosetta
directory=/app/rosetta/
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stopsignal=SIGTERM
priority=100
